# com.salestax.report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.salestax.report</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/reedhansen/Sales Tax Report/start_scheduler.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/reedhansen/Sales Tax Report/scheduler_error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/reedhansen/Sales Tax Report/scheduler_output.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/reedhansen/Sales Tax Report</string>
</dict>
</plist> 